# launchd/com.claude-memory-extract.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-memory-extract</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>REPLACE_WITH_ABSOLUTE_PATH/scripts/extract_memory.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/claude-memory-extract.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-memory-extract.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>REPLACE_WITH_HOME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PGUSER</key>
        <string>REPLACE_WITH_YOUR_USER</string>
    </dict>
</dict>
</plist>
